# scripts/com.colalin.horizon-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.colalin.horizon-daily</string>
    <key>ProgramArguments</key>
    <array>
        <!-- caffeinate keeps the Mac awake for the ~2h run after the
             04:55 pmset scheduled wake (sudo pmset repeat wakeorpoweron
             MTWRFSU 04:55:00); it releases the assertion when done. -->
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>-s</string>
        <string>/bin/bash</string>
        <string>/Users/colalin/repos/Horizon/scripts/daily-run-local.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/colalin/Library/Logs/horizon-daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/colalin/Library/Logs/horizon-daily.log</string>
</dict>
</plist>
